# deploy/roles/certbot/tasks/templates/aria-at-app-cert.timer.j2
[Unit]
Description=Attempt to renew SSL certificate once every day

[Timer]
OnBootSec=0 seconds
OnUnitActiveSec=1 day

[Install]
WantedBy=timers.target